# backup/crontab
0 0 * * * root bash /backup.sh >> /var/log/cron.log 2>&1
